# gpe-core/gpe-core.service
[Unit]
Description=GPE-Core Schwarzes-Loch-System
Documentation=file:///root/gpe-core/nordstern.md
After=network.target ollama.service
Wants=ollama.service

[Service]
# Type=forking: start_all.sh backgrounded Prozesse und beendet sich → systemd verfolgt die Kinder
Type=forking
WorkingDirectory=/root/gpe-core

# Startup via kontrolliertem Wrapper
ExecStart=/bin/bash /root/gpe-core/start_all.sh start
ExecStop=/bin/bash /root/gpe-core/start_all.sh stop
ExecReload=/bin/bash /root/gpe-core/start_all.sh restart

# Restart nur bei echtem Crash (nicht bei manuellem Stop)
Restart=on-failure
RestartSec=15
StartLimitIntervalSec=120
StartLimitBurst=3

# Umgebungsvariablen
User=root
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/root/gpe-core:/root/gpe-core/analyzer:/root/gpe-core/dispatcher"
Environment="OLLAMA_URL=http://localhost:11434"
Environment="ENABLE_SELF_IMPROVEMENT=0"

# Logging direkt zu journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gpe-core

# Resource Limits — verhindert RAM-Erschöpfung und Swap-Eskalation
# 40GB RAM-Limit für den gesamten GPE-Core Prozessbaum
MemoryMax=40G
MemorySwapMax=2G
# CPU nicht limitieren — vollen Ryzen 7950X3D nutzen
TasksMax=256

# Sicherheit
NoNewPrivileges=false
ProtectSystem=false

[Install]
WantedBy=multi-user.target
